# devonn/devonn-ecosystem.service
[Unit]
Description=DEVONN.ai Ecosystem Connector
Documentation=https://github.com/DecawDevonn/openclaw-workspace
After=network.target devonn-scheduler.service
Requires=devonn-scheduler.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/.openclaw/workspace/devonn
EnvironmentFile=/home/ubuntu/.openclaw/workspace/.env
ExecStart=/usr/bin/node ecosystem-standalone.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=devonn-ecosystem

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false

[Install]
WantedBy=multi-user.target
