# devonn/devonn-scheduler.service
[Unit]
Description=DEVONN.ai Autonomous Scheduler
Documentation=https://github.com/DecawDevonn/openclaw-workspace
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/.openclaw/workspace/devonn
EnvironmentFile=/home/ubuntu/.openclaw/workspace/.env
ExecStart=/usr/bin/node bootstrap.js
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=devonn-scheduler

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/ubuntu/.openclaw/workspace/devonn/data
ReadWritePaths=/home/ubuntu/.openclaw/workspace/devonn/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

[Install]
WantedBy=multi-user.target
